# build/content-tpl/etc/init/xbian-handle-cursor.conf
start on splash ACTION=start
stop on splash ACTION=stop

env NAME=cursor_blink

pre-start script
    mountpoint -q /sys || { stop; exit 0; }
    echo 0 > $(find /sys/devices/virtual/ -name $NAME)
end script

post-stop script
    echo 1 > $(find /sys/devices/virtual/ -name $NAME)
end script
